# launchagents/com.productivity.macbook-export.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.productivity.macbook-export</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/ryojiokuda/dev/productivity/.venv/bin/python3</string>
        <string>/Users/ryojiokuda/dev/productivity/scripts/macbook_export.py</string>
    </array>
    <key>StartInterval</key>
    <integer>60</integer>
    <key>StandardOutPath</key>
    <string>/tmp/productivity-macbook-export.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/productivity-macbook-export.log</string>
    <key>RunAtLoad</key>
    <true/>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
